# deploy/bingo.service
[Unit]
Description=Conference Bingo (Gunicorn)
After=network.target

[Service]
User=bingo
Group=bingo
WorkingDirectory=/srv/conference-bingo
EnvironmentFile=/etc/bingo.env
ExecStart=/srv/conference-bingo/venv/bin/gunicorn \
    --workers 3 \
    --bind 127.0.0.1:5000 \
    --access-logfile /var/log/bingo/access.log \
    --error-logfile /var/log/bingo/error.log \
    wsgi:app
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
